# files/etc/systemd/system/dnf-automatic.service.d/non-metered.conf
[Service]
ExecCondition=/usr/local/bin/non_metered_connection.bash
