# playbooks/files/etc/systemd/user/user-updates.timer
[Unit]
Description=User updates timer

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
