# playbooks/files/etc/systemd/system/dnf-automatic.service.d/non-metered.conf
[Service]
ExecCondition=non-metered-connection.bash
